# config/redis-fortune.conf
description "redis server for fortunes"
start on runlevel [2345];
stop on runlevel [!2345];
respawn
chdir /home/jeff/apps/redis-2.2.7/bin
exec /home/jeff/apps/redis-2.2.7/bin/redis-server